# profile-buildhost/templates/nexus.service.erb
[Unit]
Description=Nexus
After=network.target network-online.target
Wants=network-online.target

[Service]
EnvironmentFile=/opt/sonatype-nexus/conf/nexus.env
ExecStart=/opt/sonatype-nexus/bin/nexus start
ExecStop=/opt/sonatype-nexus/bin/nexus stop
ExecReload=/opt/sonatype-nexus/bin/nexus restart
PIDFile=/opt/sonatype-nexus/bin/jsw/linux-x86-64/nexus.pid
Type=forking
User=nexus


[Install]
WantedBy=multi-user.target
